# docs/RHEL_Centos/uwsgi.conf
description "uWSGI Emperor"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /var/www/datamart/env/bin/uwsgi --emperor /etc/uwsgi/vassals/
